# automation/msm-reporter.service
[Unit]
Description=MSM terminal Markdown reporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nnv
Group=nnv
WorkingDirectory=/home/nnv/MSM-Research-Lab
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/python3 -B /usr/local/lib/msm-orchestrator/msm_reporter.py --repo /home/nnv/MSM-Research-Lab --state-dir /home/nnv/.local/state/msm-orchestrator --poll 10
Restart=on-failure
RestartSec=15
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/home/nnv/.local/state/msm-orchestrator

[Install]
WantedBy=multi-user.target
