# automation/msm-dashboard.service
[Unit]
Description=MSM Research Lab web dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nnv
Group=nnv
ExecStart=/usr/bin/python3 -B /usr/local/lib/msm-orchestrator/msm_dashboard.py
Restart=always
RestartSec=5
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadOnlyPaths=/home/nnv/MSM-Research-Lab
ReadOnlyPaths=/home/nnv/.local/state/msm-orchestrator
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target