# automation/runner.timer
[Unit]
Description=Run MSM Codex task runner every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
RandomizedDelaySec=20s
Unit=msm-codex-runner.service

[Install]
WantedBy=timers.target
