# automation/runner.timer
[Unit]
Description=Run MSM Codex task runner every minute

[Timer]
OnActiveSec=15s
OnUnitInactiveSec=1min
RandomizedDelaySec=5s
AccuracySec=5s
Unit=msm-codex-runner.service

[Install]
WantedBy=timers.target
