# automation/runner.timer
[Unit]
Description=Run MSM Codex task runner every five minutes

[Timer]
OnActiveSec=30s
OnUnitInactiveSec=5min
RandomizedDelaySec=20s
AccuracySec=10s
Unit=msm-codex-runner.service

[Install]
WantedBy=timers.target
